# echo-attack-dashboard/data/G0026.json
{
  "attack_id": "G0026",
  "notes": null,
  "created": "31 May 2017",
  "modified": "11 April 2024",
  "name": "APT18",
  "version": "2.2",
  "descr": "[APT18](https://attack.mitre.org/groups/G0026) is a threat group that has operated since at least 2009 and has targeted a range of industries, including technology, manufacturing, human rights groups, government, and medical. (Citation: Dell Lateral Movement)",
  "technique_table_data": [
    {
      "domain": "enterprise",
      "id": "T1071",
      "name": "Application Layer Protocol",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "001",
          "name": "Web Protocols",
          "descr": "[APT18](https://attack.mitre.org/groups/G0026) uses HTTP for C2 communications.(Citation: PaloAlto DNS Requests May 2016)",
          "subtechniques": []
        },
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "004",
          "name": "DNS",
          "descr": "[APT18](https://attack.mitre.org/groups/G0026) uses DNS for C2 communications.(Citation: PaloAlto DNS Requests May 2016)",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1547",
      "name": "Boot or Logon Autostart Execution",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "001",
          "name": "Registry Run Keys / Startup Folder",
          "descr": "[APT18](https://attack.mitre.org/groups/G0026) establishes persistence via the <code>HKCU\\Software\\Microsoft\\Windows\\CurrentVersion\\Run</code> key.(Citation: Anomali Evasive Maneuvers July 2015)(Citation: PaloAlto DNS Requests May 2016)",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1059",
      "name": "Command and Scripting Interpreter",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "003",
          "name": "Windows Command Shell",
          "descr": "[APT18](https://attack.mitre.org/groups/G0026) uses cmd.exe to execute commands on the victim’s machine.(Citation: PaloAlto DNS Requests May 2016)(Citation: Anomali Evasive Maneuvers July 2015)",
          "subtechniques": []
        }
      ]
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1133",
      "name": "External Remote Services",
      "descr": "[APT18](https://attack.mitre.org/groups/G0026) actors leverage legitimate credentials to log into external remote services.(Citation: RSA2017 Detect and Respond Adair)",
      "subtechniques": []
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1083",
      "name": "File and Directory Discovery",
      "descr": "[APT18](https://attack.mitre.org/groups/G0026) can list files information for specific directories.(Citation: PaloAlto DNS Requests May 2016)",
      "subtechniques": []
    },
    {
      "domain": "enterprise",
      "id": "T1070",
      "name": "Indicator Removal",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "004",
          "name": "File Deletion",
          "descr": "[APT18](https://attack.mitre.org/groups/G0026) actors deleted tools and batch files from victim systems.(Citation: Dell Lateral Movement)",
          "subtechniques": []
        }
      ]
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1105",
      "name": "Ingress Tool Transfer",
      "descr": "[APT18](https://attack.mitre.org/groups/G0026) can upload a file to the victim’s machine.(Citation: PaloAlto DNS Requests May 2016)",
      "subtechniques": []
    },
    {
      "domain": "enterprise",
      "id": "T1027",
      "name": "Obfuscated Files or Information",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "013",
          "name": "Encrypted/Encoded File",
          "descr": "[APT18](https://attack.mitre.org/groups/G0026) obfuscates strings in the payload.(Citation: PaloAlto DNS Requests May 2016)",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1053",
      "name": "Scheduled Task/Job",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "002",
          "name": "At",
          "descr": "[APT18](https://attack.mitre.org/groups/G0026) actors used the native [at](https://attack.mitre.org/software/S0110) Windows task scheduler tool to use scheduled tasks for execution on a victim network.(Citation: Dell Lateral Movement)",
          "subtechniques": []
        }
      ]
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1082",
      "name": "System Information Discovery",
      "descr": "[APT18](https://attack.mitre.org/groups/G0026) can collect system information from the victim’s machine.(Citation: PaloAlto DNS Requests May 2016)",
      "subtechniques": []
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1078",
      "name": "Valid Accounts",
      "descr": "[APT18](https://attack.mitre.org/groups/G0026) actors leverage legitimate credentials to log into external remote services.(Citation: RSA2017 Detect and Respond Adair)",
      "subtechniques": []
    }
  ],
  "campaign_data": [],
  "add_campaign_ref": false,
  "software_data": [
    {
      "id": "S0106",
      "name": "cmd",
      "descr": "(Citation: Dell Lateral Movement)",
      "techniques": [
        {
          "parent_id": "T1059",
          "id": "003",
          "name": "Command and Scripting Interpreter",
          "sub_name": "Windows Command Shell"
        },
        {
          "id": "T1083",
          "name": "File and Directory Discovery"
        },
        {
          "parent_id": "T1070",
          "id": "004",
          "name": "Indicator Removal",
          "sub_name": "File Deletion"
        },
        {
          "id": "T1105",
          "name": "Ingress Tool Transfer"
        },
        {
          "id": "T1570",
          "name": "Lateral Tool Transfer"
        },
        {
          "id": "T1082",
          "name": "System Information Discovery"
        }
      ]
    },
    {
      "id": "S0032",
      "name": "gh0st RAT",
      "descr": "(Citation: RSA2017 Detect and Respond Adair)",
      "techniques": [
        {
          "parent_id": "T1547",
          "id": "001",
          "name": "Boot or Logon Autostart Execution",
          "sub_name": "Registry Run Keys / Startup Folder"
        },
        {
          "id": "T1059",
          "name": "Command and Scripting Interpreter"
        },
        {
          "parent_id": "T1543",
          "id": "003",
          "name": "Create or Modify System Process",
          "sub_name": "Windows Service"
        },
        {
          "parent_id": "T1132",
          "id": "001",
          "name": "Data Encoding",
          "sub_name": "Standard Encoding"
        },
        {
          "id": "T1140",
          "name": "Deobfuscate/Decode Files or Information"
        },
        {
          "parent_id": "T1568",
          "id": "001",
          "name": "Dynamic Resolution",
          "sub_name": "Fast Flux DNS"
        },
        {
          "parent_id": "T1573",
          "id": "001",
          "name": "Encrypted Channel",
          "sub_name": "Symmetric Cryptography"
        },
        {
          "id": "T1573",
          "name": "Encrypted Channel"
        },
        {
          "parent_id": "T1574",
          "id": "001",
          "name": "Hijack Execution Flow",
          "sub_name": "DLL"
        },
        {
          "parent_id": "T1070",
          "id": "001",
          "name": "Indicator Removal",
          "sub_name": "Clear Windows Event Logs"
        },
        {
          "parent_id": "T1070",
          "id": "004",
          "name": "Indicator Removal",
          "sub_name": "File Deletion"
        },
        {
          "id": "T1105",
          "name": "Ingress Tool Transfer"
        },
        {
          "parent_id": "T1056",
          "id": "001",
          "name": "Input Capture",
          "sub_name": "Keylogging"
        },
        {
          "id": "T1112",
          "name": "Modify Registry"
        },
        {
          "id": "T1106",
          "name": "Native API"
        },
        {
          "id": "T1095",
          "name": "Non-Application Layer Protocol"
        },
        {
          "id": "T1057",
          "name": "Process Discovery"
        },
        {
          "id": "T1055",
          "name": "Process Injection"
        },
        {
          "id": "T1012",
          "name": "Query Registry"
        },
        {
          "id": "T1113",
          "name": "Screen Capture"
        },
        {
          "id": "T1129",
          "name": "Shared Modules"
        },
        {
          "parent_id": "T1218",
          "id": "011",
          "name": "System Binary Proxy Execution",
          "sub_name": "Rundll32"
        },
        {
          "id": "T1082",
          "name": "System Information Discovery"
        },
        {
          "parent_id": "T1569",
          "id": "002",
          "name": "System Services",
          "sub_name": "Service Execution"
        }
      ]
    },
    {
      "id": "S0071",
      "name": "hcdLoader",
      "descr": "(Citation: Dell Lateral Movement)(Citation: ThreatStream Evasion Analysis)",
      "techniques": [
        {
          "parent_id": "T1059",
          "id": "003",
          "name": "Command and Scripting Interpreter",
          "sub_name": "Windows Command Shell"
        },
        {
          "parent_id": "T1543",
          "id": "003",
          "name": "Create or Modify System Process",
          "sub_name": "Windows Service"
        }
      ]
    },
    {
      "id": "S0070",
      "name": "HTTPBrowser",
      "descr": "(Citation: RSA2017 Detect and Respond Adair)",
      "techniques": [
        {
          "parent_id": "T1071",
          "id": "004",
          "name": "Application Layer Protocol",
          "sub_name": "DNS"
        },
        {
          "parent_id": "T1071",
          "id": "001",
          "name": "Application Layer Protocol",
          "sub_name": "Web Protocols"
        },
        {
          "parent_id": "T1547",
          "id": "001",
          "name": "Boot or Logon Autostart Execution",
          "sub_name": "Registry Run Keys / Startup Folder"
        },
        {
          "parent_id": "T1059",
          "id": "003",
          "name": "Command and Scripting Interpreter",
          "sub_name": "Windows Command Shell"
        },
        {
          "id": "T1083",
          "name": "File and Directory Discovery"
        },
        {
          "parent_id": "T1574",
          "id": "001",
          "name": "Hijack Execution Flow",
          "sub_name": "DLL"
        },
        {
          "parent_id": "T1070",
          "id": "004",
          "name": "Indicator Removal",
          "sub_name": "File Deletion"
        },
        {
          "id": "T1105",
          "name": "Ingress Tool Transfer"
        },
        {
          "parent_id": "T1056",
          "id": "001",
          "name": "Input Capture",
          "sub_name": "Keylogging"
        },
        {
          "parent_id": "T1036",
          "id": "005",
          "name": "Masquerading",
          "sub_name": "Match Legitimate Resource Name or Location"
        },
        {
          "id": "T1027",
          "name": "Obfuscated Files or Information"
        }
      ]
    },
    {
      "id": "S0124",
      "name": "Pisloader",
      "descr": "(Citation: Palo Alto DNS Requests)",
      "techniques": [
        {
          "parent_id": "T1071",
          "id": "004",
          "name": "Application Layer Protocol",
          "sub_name": "DNS"
        },
        {
          "parent_id": "T1547",
          "id": "001",
          "name": "Boot or Logon Autostart Execution",
          "sub_name": "Registry Run Keys / Startup Folder"
        },
        {
          "parent_id": "T1059",
          "id": "003",
          "name": "Command and Scripting Interpreter",
          "sub_name": "Windows Command Shell"
        },
        {
          "parent_id": "T1132",
          "id": "001",
          "name": "Data Encoding",
          "sub_name": "Standard Encoding"
        },
        {
          "id": "T1083",
          "name": "File and Directory Discovery"
        },
        {
          "id": "T1105",
          "name": "Ingress Tool Transfer"
        },
        {
          "id": "T1027",
          "name": "Obfuscated Files or Information"
        },
        {
          "id": "T1082",
          "name": "System Information Discovery"
        },
        {
          "id": "T1016",
          "name": "System Network Configuration Discovery"
        }
      ]
    }
  ],
  "add_software_ref": true,
  "alias_descriptions": [
    {
      "name": "TG-0416",
      "descr": "(Citation: ThreatStream Evasion Analysis)(Citation: Anomali Evasive Maneuvers July 2015)"
    },
    {
      "name": "Dynamite Panda",
      "descr": "(Citation: ThreatStream Evasion Analysis)(Citation: Anomali Evasive Maneuvers July 2015)"
    },
    {
      "name": "Threat Group-0416",
      "descr": "(Citation: ThreatStream Evasion Analysis)"
    }
  ],
  "citations": {
    "current_number": 0,
    "Dell Lateral Movement": {
      "description": "Carvey, H.. (2014, September 2). Where you AT?: Indicators of lateral movement using at.exe on Windows 7 systems. Retrieved January 25, 2016.",
      "number": null,
      "url": "http://www.secureworks.com/resources/blog/where-you-at-indicators-of-lateral-movement-using-at-exe-on-windows-7-systems/"
    },
    "Anomali Evasive Maneuvers July 2015": {
      "description": "Shelmire, A. (2015, July 06). Evasive Maneuvers by the Wekby group with custom ROP-packing and DNS covert channels. Retrieved November 15, 2018.",
      "number": null,
      "url": "https://www.anomali.com/blog/evasive-maneuvers-the-wekby-group-attempts-to-evade-analysis-via-custom-rop"
    },
    "ThreatStream Evasion Analysis": {
      "description": "Shelmire, A.. (2015, July 6). Evasive Maneuvers. Retrieved January 22, 2016.",
      "number": null,
      "url": "https://www.threatstream.com/blog/evasive-maneuvers-the-wekby-group-attempts-to-evade-analysis-via-custom-rop"
    },
    "RSA2017 Detect and Respond Adair": {
      "description": "Adair, S. (2017, February 17). Detecting and Responding to Advanced Threats within Exchange Environments. Retrieved November 17, 2024.",
      "number": null,
      "url": "https://web.archive.org/web/20210803040540/https://published-prd.lanyonevents.com/published/rsaus17/sessionsFiles/5009/HTA-F02-Detecting-and-Responding-to-Advanced-Threats-within-Exchange-Environments.pdf"
    },
    "PaloAlto DNS Requests May 2016": {
      "description": "Grunzweig, J., et al. (2016, May 24). New Wekby Attacks Use DNS Requests As Command and Control Mechanism. Retrieved November 15, 2018.",
      "number": null,
      "url": "https://researchcenter.paloaltonetworks.com/2016/05/unit42-new-wekby-attacks-use-dns-requests-as-command-and-control-mechanism/"
    },
    "Palo Alto DNS Requests": {
      "description": "Grunzweig, J., et al. (2016, May 24). New Wekby Attacks Use DNS Requests As Command and Control Mechanism. Retrieved August 17, 2016.",
      "number": null,
      "url": "http://researchcenter.paloaltonetworks.com/2016/05/unit42-new-wekby-attacks-use-dns-requests-as-command-and-control-mechanism/"
    }
  },
  "aliases_list": [
    "TG-0416",
    "Dynamite Panda",
    "Threat Group-0416"
  ]
}