# echo-attack-dashboard/data/G0099.json
{
  "attack_id": "G0099",
  "notes": null,
  "created": "05 May 2020",
  "modified": "25 April 2025",
  "name": "APT-C-36",
  "version": "1.1",
  "contributors_list": [
    "Jose Luis Sánchez Martinez"
  ],
  "descr": "[APT-C-36](https://attack.mitre.org/groups/G0099) is a suspected South America espionage group that has been active since at least 2018. The group mainly targets Colombian government institutions as well as important corporations in the financial sector, petroleum industry, and professional manufacturing.(Citation: QiAnXin APT-C-36 Feb2019)",
  "technique_table_data": [
    {
      "domain": "enterprise",
      "id": "T1059",
      "name": "Command and Scripting Interpreter",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "005",
          "name": "Visual Basic",
          "descr": "[APT-C-36](https://attack.mitre.org/groups/G0099) has embedded a VBScript within a malicious Word document which is executed upon the document opening.(Citation: QiAnXin APT-C-36 Feb2019)",
          "subtechniques": []
        }
      ]
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1105",
      "name": "Ingress Tool Transfer",
      "descr": "[APT-C-36](https://attack.mitre.org/groups/G0099) has downloaded binary data from a specified domain after the malicious document is opened.(Citation: QiAnXin APT-C-36 Feb2019)",
      "subtechniques": []
    },
    {
      "domain": "enterprise",
      "id": "T1036",
      "name": "Masquerading",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "004",
          "name": "Masquerade Task or Service",
          "descr": "[APT-C-36](https://attack.mitre.org/groups/G0099) has disguised its scheduled tasks as those used by Google.(Citation: QiAnXin APT-C-36 Feb2019)",
          "subtechniques": []
        }
      ]
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1571",
      "name": "Non-Standard Port",
      "descr": "[APT-C-36](https://attack.mitre.org/groups/G0099) has used port 4050 for C2 communications.(Citation: QiAnXin APT-C-36 Feb2019)",
      "subtechniques": []
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1027",
      "name": "Obfuscated Files or Information",
      "descr": "[APT-C-36](https://attack.mitre.org/groups/G0099) has used ConfuserEx to obfuscate its variant of [Imminent Monitor](https://attack.mitre.org/software/S0434), compressed payload and RAT packages, and password protected encrypted email attachments to avoid detection.(Citation: QiAnXin APT-C-36 Feb2019)",
      "subtechniques": []
    },
    {
      "domain": "enterprise",
      "id": "T1588",
      "name": "Obtain Capabilities",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "002",
          "name": "Tool",
          "descr": "[APT-C-36](https://attack.mitre.org/groups/G0099) obtained and used a modified variant of [Imminent Monitor](https://attack.mitre.org/software/S0434).(Citation: QiAnXin APT-C-36 Feb2019)",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1566",
      "name": "Phishing",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "001",
          "name": "Spearphishing Attachment",
          "descr": "[APT-C-36](https://attack.mitre.org/groups/G0099) has used spearphishing emails with password protected RAR attachment to avoid being detected by the email gateway.(Citation: QiAnXin APT-C-36 Feb2019) ",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1053",
      "name": "Scheduled Task/Job",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "005",
          "name": "Scheduled Task",
          "descr": "[APT-C-36](https://attack.mitre.org/groups/G0099) has used a macro function to set scheduled tasks, disguised as those used by Google.(Citation: QiAnXin APT-C-36 Feb2019)",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1204",
      "name": "User Execution",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "002",
          "name": "Malicious File",
          "descr": "[APT-C-36](https://attack.mitre.org/groups/G0099) has prompted victims to accept macros in order to execute the subsequent payload.(Citation: QiAnXin APT-C-36 Feb2019)",
          "subtechniques": []
        }
      ]
    }
  ],
  "campaign_data": [],
  "add_campaign_ref": false,
  "software_data": [
    {
      "id": "S0434",
      "name": "Imminent Monitor",
      "descr": "(Citation: QiAnXin APT-C-36 Feb2019)",
      "techniques": [
        {
          "id": "T1123",
          "name": "Audio Capture"
        },
        {
          "id": "T1059",
          "name": "Command and Scripting Interpreter"
        },
        {
          "parent_id": "T1555",
          "id": "003",
          "name": "Credentials from Password Stores",
          "sub_name": "Credentials from Web Browsers"
        },
        {
          "id": "T1140",
          "name": "Deobfuscate/Decode Files or Information"
        },
        {
          "id": "T1041",
          "name": "Exfiltration Over C2 Channel"
        },
        {
          "id": "T1083",
          "name": "File and Directory Discovery"
        },
        {
          "parent_id": "T1564",
          "id": "001",
          "name": "Hide Artifacts",
          "sub_name": "Hidden Files and Directories"
        },
        {
          "parent_id": "T1562",
          "id": "001",
          "name": "Impair Defenses",
          "sub_name": "Disable or Modify Tools"
        },
        {
          "parent_id": "T1070",
          "id": "004",
          "name": "Indicator Removal",
          "sub_name": "File Deletion"
        },
        {
          "parent_id": "T1056",
          "id": "001",
          "name": "Input Capture",
          "sub_name": "Keylogging"
        },
        {
          "id": "T1106",
          "name": "Native API"
        },
        {
          "id": "T1027",
          "name": "Obfuscated Files or Information"
        },
        {
          "id": "T1057",
          "name": "Process Discovery"
        },
        {
          "parent_id": "T1021",
          "id": "001",
          "name": "Remote Services",
          "sub_name": "Remote Desktop Protocol"
        },
        {
          "parent_id": "T1496",
          "id": "001",
          "name": "Resource Hijacking",
          "sub_name": "Compute Hijacking"
        },
        {
          "id": "T1125",
          "name": "Video Capture"
        }
      ]
    }
  ],
  "add_software_ref": true,
  "alias_descriptions": [
    {
      "name": "Blind Eagle",
      "descr": "(Citation: QiAnXin APT-C-36 Feb2019)"
    }
  ],
  "citations": {
    "current_number": 0,
    "QiAnXin APT-C-36 Feb2019": {
      "description": "QiAnXin Threat Intelligence Center. (2019, February 18). APT-C-36: Continuous Attacks Targeting Colombian Government Institutions and Corporations. Retrieved May 5, 2020.",
      "number": null,
      "url": "https://web.archive.org/web/20190625182633if_/https://ti.360.net/blog/articles/apt-c-36-continuous-attacks-targeting-colombian-government-institutions-and-corporations-en/"
    }
  },
  "aliases_list": [
    "Blind Eagle"
  ]
}