# echo-attack-dashboard/data/G0106.json
{
  "attack_id": "G0106",
  "notes": null,
  "created": "26 May 2020",
  "modified": "25 April 2025",
  "name": "Rocke",
  "version": "1.0",
  "descr": "[Rocke](https://attack.mitre.org/groups/G0106) is an alleged Chinese-speaking adversary whose primary objective appeared to be cryptojacking, or stealing victim system resources for the purposes of mining cryptocurrency. The name [Rocke](https://attack.mitre.org/groups/G0106) comes from the email address \"rocke@live.cn\" used to create the wallet which held collected cryptocurrency. Researchers have detected overlaps between [Rocke](https://attack.mitre.org/groups/G0106) and the Iron Cybercrime Group, though this attribution has not been confirmed.(Citation: Talos Rocke August 2018)",
  "technique_table_data": [
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1071",
      "name": "Application Layer Protocol",
      "descr": "[Rocke](https://attack.mitre.org/groups/G0106) issued wget requests from infected systems to the C2.(Citation: Talos Rocke August 2018)",
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "001",
          "name": "Web Protocols",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has executed wget and curl commands to Pastebin over the HTTPS protocol.(Citation: Anomali Rocke March 2019)",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1547",
      "name": "Boot or Logon Autostart Execution",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "001",
          "name": "Registry Run Keys / Startup Folder",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106)'s miner has created UPX-packed files in the Windows Start Menu Folder.(Citation: Talos Rocke August 2018)\t",
          "subtechniques": []
        }
      ]
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1037",
      "name": "Boot or Logon Initialization Scripts",
      "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has installed an \"init.d\" startup script to maintain persistence.(Citation: Anomali Rocke March 2019)\t",
      "subtechniques": []
    },
    {
      "domain": "enterprise",
      "id": "T1059",
      "name": "Command and Scripting Interpreter",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "004",
          "name": "Unix Shell",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) used shell scripts to run commands which would obtain persistence and execute the cryptocurrency mining malware.(Citation: Talos Rocke August 2018)\t",
          "subtechniques": []
        },
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "006",
          "name": "Python",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has used Python-based malware to install and spread their coinminer.(Citation: Anomali Rocke March 2019)",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1543",
      "name": "Create or Modify System Process",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "002",
          "name": "Systemd Service",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has installed a systemd service script to maintain persistence.(Citation: Anomali Rocke March 2019)\t",
          "subtechniques": []
        }
      ]
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1140",
      "name": "Deobfuscate/Decode Files or Information",
      "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has extracted tar.gz files after downloading them from a C2 server.(Citation: Talos Rocke August 2018)",
      "subtechniques": []
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1190",
      "name": "Exploit Public-Facing Application",
      "descr": "[Rocke](https://attack.mitre.org/groups/G0106) exploited Apache Struts, Oracle WebLogic (CVE-2017-10271), and Adobe ColdFusion (CVE-2017-3066) vulnerabilities to deliver malware.(Citation: Talos Rocke August 2018)(Citation: Unit 42 Rocke January 2019)",
      "subtechniques": []
    },
    {
      "domain": "enterprise",
      "id": "T1222",
      "name": "File and Directory Permissions Modification",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "002",
          "name": "Linux and Mac File and Directory Permissions Modification",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has changed file permissions of files so they could not be modified.(Citation: Anomali Rocke March 2019)\t",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1564",
      "name": "Hide Artifacts",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "001",
          "name": "Hidden Files and Directories",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) downloaded a file \"libprocesshider\", which could hide files on the target system.(Citation: Talos Rocke August 2018)(Citation: Unit 42 Rocke January 2019)",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1574",
      "name": "Hijack Execution Flow",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "006",
          "name": "Dynamic Linker Hijacking",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has modified /etc/ld.so.preload to hook libc functions in order to hide the installed dropper and mining software in process lists.(Citation: Anomali Rocke March 2019)\t",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1562",
      "name": "Impair Defenses",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "001",
          "name": "Disable or Modify Tools",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) used scripts which detected and uninstalled antivirus software.(Citation: Talos Rocke August 2018)(Citation: Unit 42 Rocke January 2019)",
          "subtechniques": []
        },
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "004",
          "name": "Disable or Modify System Firewall",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) used scripts which killed processes and added firewall rules to block traffic related to other cryptominers.(Citation: Talos Rocke August 2018)\t",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1070",
      "name": "Indicator Removal",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "002",
          "name": "Clear Linux or Mac System Logs",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has cleared log files within the /var/log/ folder.(Citation: Anomali Rocke March 2019)",
          "subtechniques": []
        },
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "004",
          "name": "File Deletion",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has deleted files on infected machines.(Citation: Anomali Rocke March 2019)\t",
          "subtechniques": []
        },
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "006",
          "name": "Timestomp",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has changed the time stamp of certain files.(Citation: Anomali Rocke March 2019)\t",
          "subtechniques": []
        }
      ]
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1105",
      "name": "Ingress Tool Transfer",
      "descr": "[Rocke](https://attack.mitre.org/groups/G0106) used malware to download additional malicious files to the target system.(Citation: Talos Rocke August 2018)\t",
      "subtechniques": []
    },
    {
      "domain": "enterprise",
      "id": "T1036",
      "name": "Masquerading",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "005",
          "name": "Match Legitimate Resource Name or Location",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has used shell scripts which download mining executables and saves them with the filename \"java\".(Citation: Talos Rocke August 2018)",
          "subtechniques": []
        }
      ]
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1046",
      "name": "Network Service Discovery",
      "descr": "[Rocke](https://attack.mitre.org/groups/G0106) conducted scanning for exposed TCP port 7001 as well as SSH and Redis servers.(Citation: Talos Rocke August 2018)(Citation: Anomali Rocke March 2019)",
      "subtechniques": []
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1571",
      "name": "Non-Standard Port",
      "descr": "[Rocke](https://attack.mitre.org/groups/G0106)'s miner connects to a C2 server using port 51640.(Citation: Anomali Rocke March 2019)\t",
      "subtechniques": []
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1027",
      "name": "Obfuscated Files or Information",
      "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has modified UPX headers after packing files to break unpackers.(Citation: Anomali Rocke March 2019)",
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "002",
          "name": "Software Packing",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106)'s miner has created UPX-packed files in the Windows Start Menu Folder.(Citation: Talos Rocke August 2018)(Citation: Unit 42 Rocke January 2019)(Citation: Anomali Rocke March 2019)",
          "subtechniques": []
        },
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "004",
          "name": "Compile After Delivery",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has compiled malware, delivered to victims as .c files, with the GNU Compiler Collection (GCC).(Citation: Anomali Rocke March 2019)\t",
          "subtechniques": []
        }
      ]
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1057",
      "name": "Process Discovery",
      "descr": "[Rocke](https://attack.mitre.org/groups/G0106) can detect a running process's PID on the infected machine.(Citation: Anomali Rocke March 2019)\t",
      "subtechniques": []
    },
    {
      "domain": "enterprise",
      "id": "T1055",
      "name": "Process Injection",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "002",
          "name": "Portable Executable Injection",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106)'s miner, \"TermsHost.exe\", evaded defenses by injecting itself into Windows processes, including Notepad.exe.(Citation: Talos Rocke August 2018)\t",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1021",
      "name": "Remote Services",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "004",
          "name": "SSH",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has spread its coinminer via SSH.(Citation: Anomali Rocke March 2019)\t",
          "subtechniques": []
        }
      ]
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1018",
      "name": "Remote System Discovery",
      "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has looked for IP addresses in the known_hosts file on the infected system and attempted to SSH into them.(Citation: Talos Rocke August 2018)\t",
      "subtechniques": []
    },
    {
      "domain": "enterprise",
      "id": "T1496",
      "name": "Resource Hijacking",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "001",
          "name": "Compute Hijacking",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has distributed cryptomining malware.(Citation: Talos Rocke August 2018)(Citation: Unit 42 Rocke January 2019)",
          "subtechniques": []
        }
      ]
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1014",
      "name": "Rootkit",
      "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has modified /etc/ld.so.preload to hook libc functions in order to hide the installed dropper and mining software in process lists.(Citation: Anomali Rocke March 2019)",
      "subtechniques": []
    },
    {
      "domain": "enterprise",
      "id": "T1053",
      "name": "Scheduled Task/Job",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "003",
          "name": "Cron",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) installed a cron job that downloaded and executed files from the C2.(Citation: Talos Rocke August 2018)(Citation: Unit 42 Rocke January 2019)(Citation: Anomali Rocke March 2019)",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1518",
      "name": "Software Discovery",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "001",
          "name": "Security Software Discovery",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) used scripts which detected and uninstalled antivirus software.(Citation: Talos Rocke August 2018)(Citation: Unit 42 Rocke January 2019)",
          "subtechniques": []
        }
      ]
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1082",
      "name": "System Information Discovery",
      "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has used uname -m to collect the name and information about the infected system's kernel.(Citation: Anomali Rocke March 2019)",
      "subtechniques": []
    },
    {
      "domain": "enterprise",
      "id": "T1552",
      "name": "Unsecured Credentials",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "004",
          "name": "Private Keys",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has used SSH private keys on the infected machine to spread its coinminer throughout a network.(Citation: Anomali Rocke March 2019)",
          "subtechniques": []
        }
      ]
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1102",
      "name": "Web Service",
      "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has used Pastebin, Gitee, and GitLab for Command and Control.(Citation: Anomali Rocke March 2019)(Citation: Talos Rocke August 2018)",
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "001",
          "name": "Dead Drop Resolver",
          "descr": "[Rocke](https://attack.mitre.org/groups/G0106) has used Pastebin to check the version of beaconing malware and redirect to another Pastebin hosting updated malware.(Citation: Anomali Rocke March 2019)",
          "subtechniques": []
        }
      ]
    }
  ],
  "campaign_data": [],
  "add_campaign_ref": false,
  "software_data": [],
  "add_software_ref": false,
  "citations": {
    "current_number": 0,
    "Talos Rocke August 2018": {
      "description": "Liebenberg, D.. (2018, August 30). Rocke: The Champion of Monero Miners. Retrieved May 26, 2020.",
      "number": null,
      "url": "https://blog.talosintelligence.com/2018/08/rocke-champion-of-monero-miners.html"
    },
    "Unit 42 Rocke January 2019": {
      "description": "Xingyu, J.. (2019, January 17). Malware Used by Rocke Group Evolves to Evade Detection by Cloud Security Products. Retrieved May 26, 2020.",
      "number": null,
      "url": "https://unit42.paloaltonetworks.com/malware-used-by-rocke-group-evolves-to-evade-detection-by-cloud-security-products/"
    },
    "Anomali Rocke March 2019": {
      "description": "Anomali Labs. (2019, March 15). Rocke Evolves Its Arsenal With a New Malware Family Written in Golang. Retrieved April 24, 2019.",
      "number": null,
      "url": "https://www.anomali.com/blog/rocke-evolves-its-arsenal-with-a-new-malware-family-written-in-golang"
    }
  },
  "aliases_list": []
}