# echo-attack-dashboard/data/G1012.json
{
  "attack_id": "G1012",
  "notes": null,
  "created": "13 January 2023",
  "modified": "02 October 2024",
  "name": "CURIUM",
  "version": "3.0",
  "contributors_list": [
    "Denise Tan",
    "Wirapong Petshagun"
  ],
  "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) is an Iranian threat group, first reported in September 2019 and active since at least July 2018, targeting IT service providers in the Middle East.(Citation: Symantec Tortoiseshell 2019) [CURIUM](https://attack.mitre.org/groups/G1012) has since invested in building relationships with potential targets via social media over a period of months to establish trust and confidence before sending malware. Security researchers note [CURIUM](https://attack.mitre.org/groups/G1012) has demonstrated great patience and persistence by chatting with potential targets daily and sending benign files to help lower their security consciousness.(Citation: Microsoft Iranian Threat Actor Trends November 2021)",
  "technique_table_data": [
    {
      "domain": "enterprise",
      "id": "T1583",
      "name": "Acquire Infrastructure",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "001",
          "name": "Domains",
          "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) created domains to facilitate strategic website compromise and credential capture activities.(Citation: PWC Yellow Liderc 2023)",
          "subtechniques": []
        },
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "003",
          "name": "Virtual Private Server",
          "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) created virtual private server instances to facilitate use of malicious domains and other items.(Citation: PWC Yellow Liderc 2023)",
          "subtechniques": []
        },
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "004",
          "name": "Server",
          "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) has created dedicated servers for command and control and exfiltration purposes.(Citation: PWC Yellow Liderc 2023)",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1059",
      "name": "Command and Scripting Interpreter",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "001",
          "name": "PowerShell",
          "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) has leveraged PowerShell scripts for initial process execution and data gathering in victim environments.(Citation: Symantec Tortoiseshell 2019)",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1584",
      "name": "Compromise Infrastructure",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "006",
          "name": "Web Services",
          "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) has compromised legitimate websites to enable strategic website compromise attacks.(Citation: PWC Yellow Liderc 2023)",
          "subtechniques": []
        }
      ]
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1005",
      "name": "Data from Local System",
      "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) has exfiltrated data from a compromised machine.(Citation: Microsoft Iranian Threat Actor Trends November 2021)",
      "subtechniques": []
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1189",
      "name": "Drive-by Compromise",
      "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) has used strategic website compromise to infect victims with malware such as [IMAPLoader](https://attack.mitre.org/software/S1152).(Citation: PWC Yellow Liderc 2023)",
      "subtechniques": []
    },
    {
      "domain": "enterprise",
      "id": "T1585",
      "name": "Establish Accounts",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "001",
          "name": "Social Media Accounts",
          "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) has established a network of fictitious social media accounts, including on Facebook and LinkedIn, to establish relationships with victims, often posing as an attractive woman.(Citation: Microsoft Iranian Threat Actor Trends November 2021)",
          "subtechniques": []
        },
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "002",
          "name": "Email Accounts",
          "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) has created dedicated email accounts for use with tools such as [IMAPLoader](https://attack.mitre.org/software/S1152).(Citation: PWC Yellow Liderc 2023)",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1048",
      "name": "Exfiltration Over Alternative Protocol",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "002",
          "name": "Exfiltration Over Asymmetric Encrypted Non-C2 Protocol",
          "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) has used SMTPS to exfiltrate collected data from victims.(Citation: PWC Yellow Liderc 2023)",
          "subtechniques": []
        }
      ]
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1041",
      "name": "Exfiltration Over C2 Channel",
      "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) has used IMAP and SMTPS for exfiltration via tools such as [IMAPLoader](https://attack.mitre.org/software/S1152).(Citation: PWC Yellow Liderc 2023)",
      "subtechniques": []
    },
    {
      "domain": "enterprise",
      "id": "T1566",
      "name": "Phishing",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "001",
          "name": "Spearphishing Attachment",
          "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) has used phishing with malicious attachments for initial access to victim environments.(Citation: PWC Yellow Liderc 2023)",
          "subtechniques": []
        },
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "003",
          "name": "Spearphishing via Service",
          "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) has used social media to deliver malicious files to victims.(Citation: Microsoft Iranian Threat Actor Trends November 2021)",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1598",
      "name": "Phishing for Information",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "003",
          "name": "Spearphishing Link",
          "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) used malicious links to adversary-controlled resources for credential harvesting.(Citation: PWC Yellow Liderc 2023)",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1505",
      "name": "Server Software Component",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "003",
          "name": "Web Shell",
          "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) has been linked to web shells following likely server compromise as an initial access vector into victim networks.(Citation: Symantec Tortoiseshell 2019)",
          "subtechniques": []
        }
      ]
    },
    {
      "domain": "enterprise",
      "id": "T1608",
      "name": "Stage Capabilities",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "004",
          "name": "Drive-by Target",
          "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) used strategic website compromise to fingerprint then target victims.(Citation: PWC Yellow Liderc 2023)",
          "subtechniques": []
        }
      ]
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1082",
      "name": "System Information Discovery",
      "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) deploys information gathering tools focused on capturing IP configuration, running application, system information, and network connectivity information.(Citation: Symantec Tortoiseshell 2019)",
      "subtechniques": []
    },
    {
      "technique_used": true,
      "domain": "enterprise",
      "id": "T1124",
      "name": "System Time Discovery",
      "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) deployed mechanisms to check system time information following strategic website compromise attacks.(Citation: PWC Yellow Liderc 2023)",
      "subtechniques": []
    },
    {
      "domain": "enterprise",
      "id": "T1204",
      "name": "User Execution",
      "technique_used": false,
      "subtechniques": [
        {
          "technique_used": true,
          "domain": "enterprise",
          "id": "002",
          "name": "Malicious File",
          "descr": "[CURIUM](https://attack.mitre.org/groups/G1012) has lured users into opening malicious files delivered via social media.(Citation: Microsoft Iranian Threat Actor Trends November 2021)",
          "subtechniques": []
        }
      ]
    }
  ],
  "campaign_data": [],
  "add_campaign_ref": false,
  "software_data": [
    {
      "id": "S1152",
      "name": "IMAPLoader",
      "descr": "[IMAPLoader](https://attack.mitre.org/software/S1152) was deployed by [CURIUM](https://attack.mitre.org/groups/G1012) as a post-exploitation payload from strategic website compromise.(Citation: PWC Yellow Liderc 2023)",
      "techniques": [
        {
          "parent_id": "T1071",
          "id": "003",
          "name": "Application Layer Protocol",
          "sub_name": "Mail Protocols"
        },
        {
          "id": "T1543",
          "name": "Create or Modify System Process"
        },
        {
          "parent_id": "T1564",
          "id": "003",
          "name": "Hide Artifacts",
          "sub_name": "Hidden Window"
        },
        {
          "parent_id": "T1574",
          "id": "014",
          "name": "Hijack Execution Flow",
          "sub_name": "AppDomainManager"
        },
        {
          "id": "T1105",
          "name": "Ingress Tool Transfer"
        },
        {
          "id": "T1106",
          "name": "Native API"
        },
        {
          "parent_id": "T1053",
          "id": "005",
          "name": "Scheduled Task/Job",
          "sub_name": "Scheduled Task"
        },
        {
          "id": "T1082",
          "name": "System Information Discovery"
        },
        {
          "id": "T1047",
          "name": "Windows Management Instrumentation"
        }
      ]
    }
  ],
  "add_software_ref": true,
  "alias_descriptions": [
    {
      "name": "Crimson Sandstorm",
      "descr": "(Citation: Microsoft Threat Actor Naming July 2023)"
    },
    {
      "name": "TA456",
      "descr": "(Citation: Microsoft Threat Actor Naming July 2023)(Citation: Proofpoint TA456 Defense Contractor July 2021)"
    },
    {
      "name": "Tortoise Shell",
      "descr": "(Citation: Microsoft Threat Actor Naming July 2023)"
    },
    {
      "name": "Yellow Liderc",
      "descr": "(Citation: PWC Yellow Liderc 2023)"
    }
  ],
  "citations": {
    "current_number": 0,
    "Microsoft Threat Actor Naming July 2023": {
      "description": "Microsoft . (2023, July 12). How Microsoft names threat actors. Retrieved November 17, 2023.",
      "number": null,
      "url": "https://learn.microsoft.com/en-us/microsoft-365/security/intelligence/microsoft-threat-actor-naming?view=o365-worldwide"
    },
    "Proofpoint TA456 Defense Contractor July 2021": {
      "description": "Miller, J. et. al. (2021, July 28). I Knew You Were Trouble: TA456 Targets Defense Contractor with Alluring Social Media Persona. Retrieved March 11, 2024.",
      "number": null,
      "url": "https://www.proofpoint.com/us/blog/threat-insight/i-knew-you-were-trouble-ta456-targets-defense-contractor-alluring-social-media"
    },
    "Microsoft Iranian Threat Actor Trends November 2021": {
      "description": "MSTIC. (2021, November 16). Evolving trends in Iranian threat actor activity – MSTIC presentation at CyberWarCon 2021. Retrieved January 12, 2023.",
      "number": null,
      "url": "https://www.microsoft.com/en-us/security/blog/2021/11/16/evolving-trends-in-iranian-threat-actor-activity-mstic-presentation-at-cyberwarcon-2021"
    },
    "PWC Yellow Liderc 2023": {
      "description": "PwC Threat Intelligence. (2023, October 25). Yellow Liderc ships its scripts and delivers IMAPLoader malware. Retrieved August 14, 2024.",
      "number": null,
      "url": "https://www.pwc.com/gx/en/issues/cybersecurity/cyber-threat-intelligence/yellow-liderc-ships-its-scripts-delivers-imaploader-malware.html"
    },
    "Symantec Tortoiseshell 2019": {
      "description": "Symantec Threat Hunter Team. (2019, September 18). Tortoiseshell Group Targets IT Providers in Saudi Arabia in Probable Supply Chain Attacks. Retrieved May 20, 2024.",
      "number": null,
      "url": "https://symantec-enterprise-blogs.security.com/blogs/threat-intelligence/tortoiseshell-apt-supply-chain"
    }
  },
  "aliases_list": [
    "Crimson Sandstorm",
    "TA456",
    "Tortoise Shell",
    "Yellow Liderc"
  ]
}